# app/crontab
* 2 * * * python3 /app/main.py >> /app/cronlog.log
